# laitshow.service
[Unit]
Description=Neopixel control server
After=network.target

[Service]
Type=simple
Restart=always
ExecStart=/home/mitchell/lights_venv/bin/python /home/mitchell/LaitShow/controller.py

[Install]
WantedBy=multi-user.target
